# ansible/roles/control-plane-alert-collector/templates/control-plane-alert-collector.timer.j2
[Unit]
Description=Run the control-plane landmine collector periodically

[Timer]
OnCalendar={{ cpa_oncalendar }}
RandomizedDelaySec={{ cpa_randomized_delay_sec }}
Persistent=true

[Install]
WantedBy=timers.target
