# tools/account_automation/deploy/crontab
# Run account automation daily at 2:00 AM (container timezone)
0 2 * * * /usr/local/bin/account-automation
